# com.onurcelik.exposure_attendance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.onurcelik.exposure_attendance</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>-u</string>
        <string>drive_monitor.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/onurcelik/Desktop/Cursor Projects/exposure attendance/exposure attendance</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/onurcelik/Desktop/Cursor Projects/exposure attendance/exposure attendance/monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/onurcelik/Desktop/Cursor Projects/exposure attendance/exposure attendance/monitor.err</string>
</dict>
</plist>
